# systemd/byrbtfunbox.timer
[Unit]
Description=ByrBT FunBox to Telegraph

[Timer]
OnBootSec=25s
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target